# config/systemd/pigpio.service
[Unit]
Description=Run the pigpio service

[Service]
ExecStart=/opt/gong/run-pigpio.sh

[Install]
WantedBy=multi-user.target

